# cd_utils/service/systemctl/deploy/expire-clients.timer
[Unit]
Description=Remover clientes expirados

[Timer]
OnActiveSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target