# cd_utils/service/systemctl/deploy/run-pending-deploys.timer
[Unit]
Description=Cd utils para deployar
After=network.target

[Timer]
OnActiveSec=10s
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target